# services/python/ai-ml-services/arcface-service/deployment/systemd/arcface-service.service
[Unit]
Description=ArcFace Face Matching Service
Documentation=https://github.com/your-org/arcface-service
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/UNIFIED_PLATFORM_COMPLETE/backend/ai-ml-services/arcface-service
Environment="PATH=/home/ubuntu/UNIFIED_PLATFORM_COMPLETE/backend/ai-ml-services/arcface-service/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/ubuntu/UNIFIED_PLATFORM_COMPLETE/backend/ai-ml-services/arcface-service/deployment/.env.staging

# Start command
ExecStart=/home/ubuntu/UNIFIED_PLATFORM_COMPLETE/backend/ai-ml-services/arcface-service/venv/bin/uvicorn \
    main:app \
    --host 0.0.0.0 \
    --port 8004 \
    --workers 4 \
    --log-level info \
    --access-log \
    --use-colors

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/arcface-service /home/ubuntu/UNIFIED_PLATFORM_COMPLETE/backend/ai-ml-services/arcface-service/cache

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arcface-service

[Install]
WantedBy=multi-user.target
